# droneaware-ble.service
[Unit]
Description=DroneAware BLE Remote ID Feeder
After=network-online.target bluetooth.target droneaware-bt-select.service
Wants=network-online.target
Requires=droneaware-bt-select.service

[Service]
Type=simple
User=root
EnvironmentFile=/opt/droneaware/config.env
Environment=PATH=/opt/droneaware/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
WorkingDirectory=/opt/droneaware
ExecStartPre=/usr/sbin/rfkill unblock bluetooth
ExecStart=/opt/droneaware/venv/bin/python3 /opt/droneaware/ble_feeder.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=droneaware-ble

[Install]
WantedBy=multi-user.target
