# droneaware-wifi.service
[Unit]
Description=DroneAware WiFi Remote ID Feeder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
EnvironmentFile=/opt/droneaware/config.env
Environment=PATH=/opt/droneaware/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
WorkingDirectory=/opt/droneaware
ExecStart=/opt/droneaware/venv/bin/python3 /opt/droneaware/wifi_feeder.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=droneaware-wifi

[Install]
WantedBy=multi-user.target
